# setup/periodic/weekly/temp-log-rotate
#!/bin/sh

mv ~/.telemetry/temperature.log ~/.telemetry/temperature-$(date +%F).log
